# onairdaemon.service
[Unit]
Description=On Air Lamp Daemon

[Service]
ExecStart=/usr/local/bin/onairdaemon
Restart=on-failure

[Install]
WantedBy=multi-user.target 
